# overlays/brash/etc/update-motd.d/10-help-text
#!/bin/bash
printf "\n"
printf " * Documentation: http://phenix-rising.srn.sandia.gov/sceptre"
printf "\n"
